# openrc/root/etc/inittab
::sysinit:openrc sysinit
::wait:openrc boot
::wait:openrc default
::ctrlaltdel:/sbin/reboot
::shutdown:openrc shutdown
::restart:openrc reboot
ttyS0::respawn:/usr/bin/env ENV=/etc/profile /bin/sh
