# openrc/root/etc/init.d/xcluster-start
#!/sbin/openrc-run

depend() {
	provide net
}

start() {
	# Don't call rcS if already called!!
	test -e /run/xcluster-start || /etc/init.d/rcS
	date >> /run/xcluster-start
}
stop() {
	return 0
}
status() {
	return 0
}
